# Update.cron
00 09 * * 3 [ $(date +\%d) -le 07 ] && ./GrabAndPush.sh
